# fs/mt7621/root_base/default/inittab
::sysinit:/default/rcS
console::respawn:/bin/sh
